# sshConfigs/chrigSystemd/wol.service
[Unit]
Description=Enable Wake On Lan

[Service]
Type=oneshot
ExecStart = /sbin/ethtool --change enp5s0 wol g

[Install]
WantedBy=basic.target
